# crontabs/dropbox.linux.crontab
*/10 * * * * $(cat ~/.dotfiles_dir_path)/scripts/dropbox_scheduler.sh
